# master-firmware/packaging/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable cvra-master
systemctl start cvra-master
