# master-firmware/packaging/prerm
#!/bin/sh

systemctl stop cvra-master
systemctl disable cvra-master
